# auto_share.timer
[Unit]
Description=Check the network mounts

[Timer]
#OnCalendar=*-*-* *:*:00
OnBootSec=5sec
OnUnitActiveSec=1min

[Install]
WantedBy=timer.target
